# build/deb/bbs-server-0.4.1/DEBIAN/postinst
#!/bin/bash
set -e

# Post-installation script for bbs-server

case "$1" in
    configure)
        echo "Configuring bbs-server..."
        
        # Create necessary directories
        install -d -m 0755 -o bbs -g bbs /var/lib/bbs
        install -d -m 0755 -o bbs -g bbs /var/lib/bbs/plugins/games
        install -d -m 0755 -o root -g root /etc/bbs
        install -d -m 0755 -o bbs -g bbs /var/log/bbs
        
        # Create default environment file if it doesn't exist
        if [ ! -f /etc/bbs/bbs.env ]; then
            cat > /etc/bbs/bbs.env << 'EOF'
# Build-a-Bot Stadium Server Configuration
# See /opt/bbs/docs/deployment/LINUX_PACKAGING_PROFILE.md for full documentation

# Packaging mode: use FHS-standard directories
BBS_PACKAGING_MODE=linux-fhs

# Dashboard admin key (set this to enable admin controls)
# BBS_DASHBOARD_ADMIN_KEY=your-secret-key-here

# Federation settings (optional)
# BBS_ENABLE_MOCK_GLOBAL_REGISTRY=true
# BBS_FEDERATION_OUTBOX_URL=https://centeralregistry.example.com/inbound
# BBS_FEDERATION_OUTBOX_TOKEN=server-api-token

# Port configuration (override in systemd service if needed)
# BBS_STADIUM_PORT=8080
# BBS_DASHBOARD_PORT=3000
EOF
            chmod 0640 /etc/bbs/bbs.env
            chown root:bbs /etc/bbs/bbs.env
            echo "Created /etc/bbs/bbs.env (default configuration)"
        fi
        
        # Enable and start the systemd service
        if command -v systemctl &>/dev/null; then
            systemctl daemon-reload
            echo ""
            echo "To start the BBS server, run:"
            echo "  sudo systemctl start bbs-server"
            echo ""
            echo "To enable auto-start on boot, run:"
            echo "  sudo systemctl enable bbs-server"
            echo ""
            echo "To view logs, run:"
            echo "  sudo journalctl -u bbs-server -f"
        fi
        ;;
esac

exit 0
